# scripts/macos/com.ra-h.sqlite-backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.ra-h.sqlite-backup</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>-lc</string>
      <string>cd /Users/elche/Code/Dev/ra-h_os && bash scripts/database/sqlite-backup.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/elche/Code/Dev/ra-h_os</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>BACKUP_DIR</key>
      <string>/Users/elche/Library/Mobile Documents/com~apple~CloudDocs/Backups/ra-h_os</string>
    </dict>

    <key>StartCalendarInterval</key>
    <dict>
      <key>Hour</key>
      <integer>2</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/elche/Library/Logs/ra-h-sqlite-backup.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/elche/Library/Logs/ra-h-sqlite-backup.log</string>
  </dict>
</plist>
